# shell/zsh/profile.d/01_basic.zsh
# cargo
. "$HOME/.cargo/env"

# Set up fzf key bindings and fuzzy completion
# require >= 0.48.0
# source <(fzf --zsh)

export TF_PLUGIN_CACHE_DIR="$HOME/.terraform.d/plugin-cache"
